# deploy/gallery.service
[Unit]
Description=Gallery app
After=network.target

[Service]
Type=simple
User=gallery
WorkingDirectory=/opt/gallery
EnvironmentFile=/etc/gallery.env
ExecStart=/opt/gallery/server/node_modules/.bin/tsx server/src/index.ts
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
